# deployment/jper_web.conf
[program:jper]
command=/home/cloo/jper/bin/gunicorn -c /home/cloo/jper/src/jper/deployment/gconf.py service.web:app
user=cloo
directory=/home/cloo/jper/src/jper
stdout_logfile=/var/log/supervisor/%(program_name)s-access.log
stderr_logfile=/var/log/supervisor/%(program_name)s-error.log
autostart=true
autorestart=true
stopasgroup=true
