# deployment/jper_scheduler.conf
[program:scheduler]
command=/home/cloo/jper/bin/python /home/cloo/jper/src/jper/service/scheduler.py
user=cloo
directory=/home/cloo/jper/src/jper
stdout_logfile=/var/log/supervisor/%(program_name)s-access.log
stderr_logfile=/var/log/supervisor/%(program_name)s-error.log
autostart=true
autorestart=true
stopasgroup=true
